# bin/freefactory-notify.service
[Unit]
Description=FreeFactory Drop Folder Monitor
After=network.target local-fs.target remote-fs.target

[Service]
Type=simple
ExecStart=/opt/FreeFactory/bin/FreeFactoryNotifyRunner.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
